# deployment/supervisor.cfg

[program:chsmysteries]
command=/home/ubuntu/charlestonmysteries/venv/bin/gunicorn -b localhost:8000 -w 4 chsmysteries:app
directory=/home/ubuntu/charlestonmysteries
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
